# debian/postinst
#!/bin/sh

update-rc.d wb-shtrixmqtt defaults
invoke-rc.d wb-shtrixmqtt start
